# scripts/podman/kolb-bot.container.in
# Kolb-Bot gateway — Podman Quadlet (rootless)
# Installed by setup-podman.sh into kolb-bot's ~/.config/containers/systemd/
# {{KOLB_BOT_HOME}} is replaced at install time.

[Unit]
Description=Kolb-Bot gateway (rootless Podman)

[Container]
Image=kolb-bot:local
ContainerName=kolb-bot
UserNS=keep-id
# Keep container UID/GID aligned with the invoking user so mounted config is readable.
User=%U:%G
Volume={{KOLB_BOT_HOME}}/.kolb-bot:/home/node/.kolb-bot:Z
EnvironmentFile={{KOLB_BOT_HOME}}/.kolb-bot/.env
Environment=HOME=/home/node
Environment=TERM=xterm-256color
PublishPort=18789:18789
PublishPort=18790:18790
Pull=never
Exec=node dist/index.js gateway --bind lan --port 18789

[Service]
TimeoutStartSec=300
Restart=on-failure

[Install]
WantedBy=default.target
